# deploy/com.claudio.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claudio.bot</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/sissifeng/.openclaw/workspace/claudio-music-agent/venv/bin/python</string>
        <string>/Users/sissifeng/.openclaw/workspace/claudio-music-agent/bot.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/sissifeng/.openclaw/workspace/claudio-music-agent</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>TELEGRAM_BOT_TOKEN</key>
        <string>${TELEGRAM_BOT_TOKEN}</string>
        <key>ELEVENLABS_API_KEY</key>
        <string>${ELEVENLABS_API_KEY}</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/sissifeng/.openclaw/workspace/claudio-music-agent/logs/claudio.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/sissifeng/.openclaw/workspace/claudio-music-agent/logs/claudio.error.log</string>
</dict>
</plist>
